# system/etc/chroot-service/profile.d/android-permissions.sh
chrootargs="$chrootargs -v /data/adb/chroot-service/android-permissions/release/aids.conf:/lib/sysusers.d/aids.conf"